# mod-tile-docker/supervisord.conf
[supervisord]
nodaemon=true

[program:renderd]
command=renderd -f -c /usr/local/etc/renderd.conf

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
